# target/linux/ramips/mt7621/base-files/etc/init.d/hc_ttyd
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
    # 启动 ttyd 和 minicom
    /usr/bin/ttyd -p 7682 --max-clients 1 minicom -D /dev/ttyS1 -b 115200 -C /mnt/sd/nxp_term.log -O timestamp=extended &
    /usr/bin/ttyd -p 7681 -d 7 /bin/login
}

stop() {
    # 停止 ttyd 服务
    killall ttyd
}

restart() {
    stop
    start
}
